# rootfs/etc/service.d/gazleolvaso/finish
#!/bin/sh
echo "[INFO] gazleolvaso service leállt"
