# rootfs/etc/service.d/gazleolvaso/run
#!/usr/bin/with-contenv bash
set -e

echo "[INFO] gazleolvaso service indul"

exec python3 /app/cli.py
